# healthcheck.cron
#!/bin/bash

set -e
set -o pipefail

log "Checking if $(basename $0) is healthy ..."

pgrep Xcron &> /dev/null

